# alpine-sync/root/etc/services.d/sync/run
#!/usr/bin/with-contenv sh

rslsync --nodaemon --log /var/log/sync.log --config /data/sync.conf