# alpine-base/root/etc/periodic/monthly/rsnapshot_monthly
#!/usr/bin/with-contenv sh

/usr/bin/rsnapshot -c /etc/rsnapshot.conf monthly
