# alpine-proxy/root/etc/services.d/consul-template/run
#!/usr/bin/with-contenv sh

exec /bin/consul-template -consul-addr=${CONSUL_ADDRESS} -config=/etc/consul-template/conf/
