# alpine-cadvisor/root/etc/services.d/cadvisor/run
#!/usr/bin/with-contenv sh

exec /usr/bin/cadvisor ${CADVISOR_OPTIONS}
