# alpine-proxy/root/etc/periodic/daily/update_cert
#!/usr/bin/with-contenv sh

/usr/bin/certbot renew  --config /etc/cli.ini --post-hook "nginx -s reload"
