# alpine-proxy/root/etc/periodic/daily/update_cert
#!/usr/bin/with-contenv sh

if curl --output /dev/null --silent --head --fail ${CERT_DOMAIN}; then
  /usr/bin/certbot renew --post-hook "nginx -s reload"
fi
