# alpine-registrator/root/etc/services.d/registrator/run
#!/usr/bin/with-contenv sh

exec /bin/registrator ${REGISTRATOR_OPTIONS} consul://${CONSUL_IP}:${CONSUL_PORT}
