# alpine-mediawiki/root/etc/services.d/memcached/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid memcache /usr/bin/memcached -v -m 128 >> /var/log/memcached.log 2>&1
